# backend/build/worker/crontab
0 * * * * /app/task.sh >> /var/log/cron.log 2>&1
